# systemd/electridrive-sync.service
[Unit]
Description=ElectriDrive Linux Safe Upload Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Edit these paths to where you cloned the repo (the folder containing app.py).
WorkingDirectory=%h/electri-city-gdrive
ExecStart=%h/electri-city-gdrive/.venv/bin/python -m electridrive.cli sync-up %h/Dokumente --remote-folder ElectriDrive/LinuxBackup
Nice=10
IOSchedulingClass=best-effort

[Install]
WantedBy=default.target
